# mail/sendmail/pkg/esendmail.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/libexec/sendmail/sendmail"

. /etc/rc.d/rc.subr

pexp="(sendmail: (accepting|rejecting).*|${daemon}.* -(q[0-9]|bd)*)"

rc_bg=YES

rc_cmd $1
